# systemd/codexcomp.service.example
# systemd user unit example (Linux/WSL). Install:
#   cp systemd/codexcomp.service.example ~/.config/systemd/user/codexcomp.service
#   systemctl --user daemon-reload && systemctl --user enable --now codexcomp
# Assumes the CLI was installed with: uv tool install codexcomp
[Unit]
Description=codexcomp: local Responses proxy folding gpt-5.5 518n-2 reasoning truncation
After=network-online.target

[Service]
ExecStart=%h/.local/bin/codexcomp
Restart=on-failure
RestartSec=2
# If your systemd user manager exports a SOCKS proxy (ALL_PROXY), httpx would
# need extra deps to honor it — drop the SOCKS vars and let HTTP(S)_PROXY win.
UnsetEnvironment=ALL_PROXY all_proxy SOCKS_PROXY socks_proxy

[Install]
WantedBy=default.target
